# deploy/systemd/sugarpy-jupyter.service
[Unit]
Description=SugarPy Jupyter Server
After=network.target

[Service]
Type=simple
User=sugarpy
WorkingDirectory=/opt/sugarpy/current
Environment=IPYTHONDIR=/opt/sugarpy/shared/.ipython
Environment=JUPYTER_CONFIG_DIR=/opt/sugarpy/shared/.jupyter
Environment=JUPYTER_RUNTIME_DIR=/opt/sugarpy/shared/.jupyter-runtime
Environment=UV_PROJECT_ENVIRONMENT=/opt/sugarpy/shared/.venv
EnvironmentFile=-/etc/sugarpy/assistant.env
Environment=SUGARPY_SECURITY_PROFILE=restricted-demo
Environment=SUGARPY_STORAGE_ROOT=/opt/sugarpy/shared/runtime
Environment=SUGARPY_ENABLE_ASSISTANT_TRACES=0
ExecStart=/opt/sugarpy/current/scripts/start-jupyter-service.sh
Restart=on-failure
RestartSec=5
MemoryMax=8G
MemorySwapMax=0
CPUQuota=200%
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ProtectKernelTunables=yes
ProtectControlGroups=yes
PrivateDevices=yes
ReadWritePaths=/opt/sugarpy/shared

[Install]
WantedBy=multi-user.target
